# deviceoswd/deb/debian/deviceos-wd.service
[Unit]
Description="Device OS watchdog reseting service"

[Service]
Restart=always
ExecStartPre=-/usr/bin/pkill deviceOSWD
ExecStart=/usr/bin/deviceOSWD -w 3000 -m 900 -s /var/lib/pelion/deviceOSkeepalive
ExecStartPost=sh -c "echo -e \"dog_stop\" | socat unix-sendto:/var/lib/pelion/deviceOSkeepalive STDIO"
ExecStop=sh -c "echo -e stop | socat unix-sendto:/var/lib/pelion/deviceOSkeepalive STDIO"

[Install]
WantedBy=basic.target
