# kubelet/deb/debian/kubelet.service
[Unit]
Description=Kubelet

[Service]
Restart=always
RestartSec=5s
ExecStart=/usr/bin/launch-kubelet.sh
