# pe-utils/deb/debian/wait-for-pelion-identity.service
[Unit]
Description=Wait for a connection to Pelion and create cridentials
Requires=edge-core.service
After=edge-core.service

[Service]
Restart=on-failure
RestartSec=5s
RemainAfterExit=true
ExecStart=/usr/bin/launch-wait-for-pelion-identity.sh

[Install]
WantedBy=multi-user.target
