# edge-proxy/deb/debian/edge-proxy.service
[Unit]
Description=Edge Proxy

[Service]
Restart=always
RestartSec=5s
ExecStart=/usr/bin/launch-edge-proxy.sh

[Install]
WantedBy=multi-user.target
